# config/telemt-compose-update.service
[Unit]
Description=Rebuild & update Telemt MTProto container (from GHCR source image)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/telemt
# Авто-обновление делегировано install-mtproto.sh: оно резолвит новейший
# тег из GHCR (для UPDATE_CHANNEL=latest), правит .env, пересобирает образ
# и пересоздаёт контейнер. Для UPDATE_CHANNEL=pinned — no-op.
ExecStart=/opt/telemt/install-mtproto.sh --auto-update
Restart=on-failure
RestartSec=60
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telemt-compose-update
